# ops/labs/systemd/hello.service
[Unit]
Description=Linux basic hello service

[Service]
Type=oneshot
ExecStart=/usr/bin/printf 'hello from systemd\n'
